# etc/profile.d/custom-path.sh
append_path "$HOME/.local/bin"
export PATH