# etc/profile.d/custom-path.sh
append_path "$HOME/.local/bin"
append_path "$HOME/go/bin"
append_path "$HOME/.tfenv/bin/"
append_path "$HOME/.cargo/bin"
export PATH
export QT_STYLE_OVERRIDE=kvantum
